# github_monitor/cli/systemd/github-monitor-indicator.service
[Unit]
Description=GitHub PR Monitor - Panel Indicator
After=github-monitor.service
Wants=github-monitor.service

[Service]
Type=simple
ExecStart=@@GITHUB_MONITOR_INDICATOR_EXEC@@
Restart=on-failure
RestartSec=10

# Security hardening (lighter than daemon — no file writes needed)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=default.target
