# github_monitor/cli/systemd/github-monitor.service
[Unit]
Description=GitHub PR Monitor
After=network-online.target dbus.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=@@GITHUB_MONITOR_EXEC@@
ExecReload=kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
Environment=GITHUB_TOKEN=

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/github-monitor

[Install]
WantedBy=default.target
